# wttr.timer
[Unit]
Description=Run wttr hourly after bootup

[Timer]
OnCalendar=hourly
Unit=motd.service

[Install]
WantedBy=basic.target
